# session-mkdir@.service
[Unit]
Description=Make session directory for %i
RefuseManualStart=true

[Service]
Type=oneshot
ExecStart=/usr/bin/mkdir -m700 -p %t/session/%i/