# session@.target
[Unit]
Description=Session %i
Wants=session-mkdir@%i.service
After=session-mkdir@%i.service
