# systemd-user-manager.service.in
[Unit]
Description=User Manager
DefaultDependencies=false
After=sockets.target dbus.target
Before=shutdown.target

[Service]
Slice=background.slice
ExecStart={PREFIX}/lib/systemd/systemd-user-manager

[Install]
WantedBy=default.target
